# cleo.service
# Template systemd service for Cleo.
# setup.sh generates the real service file with correct absolute paths.
# If installing manually, replace %h with your home directory path.
[Unit]
Description=Cleo Signal Gateway
After=network.target signal-cli.service
Requires=signal-cli.service

[Service]
Type=simple
WorkingDirectory=%h/code/cleo
ExecStart=%h/code/cleo/venv/bin/python gateway.py
Restart=on-failure
RestartSec=10
EnvironmentFile=%h/code/cleo/.env

[Install]
WantedBy=default.target
